# systemd/iasantiago-backup.service
[Unit]
Description=IASantiago backups (topics, qdrant, whoosh)

[Service]
Type=oneshot
WorkingDirectory=/opt/iasantiago-rag
ExecStart=/usr/bin/make backup
