# systemd/iasantiago-rag-eval.service
[Unit]
Description=IASantiago offline eval (Recall/MRR) and write CSV

[Service]
Type=oneshot
WorkingDirectory=/opt/iasantiago-rag
ExecStart=/bin/bash -lc 'python3 - << "PY"\nimport json,requests,csv\ncases=[\n  {\"query\":\"definición de ley de Ohm\",\"topic\":\"Electronics\",\"relevant_files\":[\"/opt/iasantiago-rag/topics/Electronics/sample1.pdf\"]},\n  {\"query\":\"ácidos y bases\",\"topic\":\"Chemistry\",\"relevant_files\":[\"/opt/iasantiago-rag/topics/Chemistry/sample1.pdf\"]}\n]\nurl=\"http://127.0.0.1:8001/v1/eval/offline\"\nr=requests.post(url,json=cases,timeout=120)\nr.raise_for_status()\nout=r.json()\nprint(out)\nwith open(\"eval_summary.csv\",\"w\",newline=\"\",encoding=\"utf-8\") as f:\n    w=csv.writer(f)\n    agg=out[\"aggregate\"]\n    w.writerow([\"n\",\"Recall@1\",\"Recall@3\",\"Recall@5\",\"Recall@10\",\"MRR\"])\n    w.writerow([agg[\"n\"],agg.get(\"Recall@1\"),agg.get(\"Recall@3\"),agg.get(\"Recall@5\"),agg.get(\"Recall@10\"),agg.get(\"MRR\")])\nPY'
