# systemd/iasantiago-rag.service
[Unit]
Description=IASantiago RAG stack (docker compose)
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
WorkingDirectory=/docker/iaburuaga-rag
ExecStart=/bin/bash -c '/usr/bin/docker compose up -d --no-start && /usr/bin/docker compose start qdrant vllm rag-api openwebui oauth2-proxy'
ExecStop=/usr/bin/docker compose stop
RemainAfterExit=yes
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
