# systemd/iasantiago-backup.service
[Unit]
Description=IASantiago backups (topics, qdrant, whoosh)

[Service]
Type=oneshot
WorkingDirectory=/docker/iaburuaga-rag
ExecStart=/usr/bin/make backup
