# systemd/iasantiago-rag-eval.timer
[Unit]
Description=Nightly offline eval at 02:30

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
Unit=iasantiago-rag-eval.service

[Install]
WantedBy=timers.target
